# docker/s6-overlay/s6-rc.d/celery-worker/run
#!/command/with-contenv sh
poetry run -C /opt celery -A src.Emailkasten worker --loglevel=${CELERY_LOG_LEVEL:-INFO} --without-gossip --without-mingle
