# docker/s6-overlay/s6-rc.d/celery-beat/run
#!/command/with-contenv sh
poetry run -C /opt celery -A src.Emailkasten beat --loglevel=${CELERY_LOG_LEVEL:-INFO} -S django
