# docker/s6-overlay/s6-rc.d/migrations/run
#!/command/with-contenv sh
python3 /opt/manage.py migrate --no-input
