# docker/s6-overlay/s6-rc.d/migrations/run
#!/command/with-contenv sh
poetry run -C /opt python manage.py migrate --no-input
