# docker/s6-overlay/s6-rc.d/createsuperuser/run
#!/command/with-contenv sh
python3 /opt/manage.py createsuperuser --username=admin --email=  --noinput
s6-svc -Od . # needed to cover up failure if user already exists
