# docker/s6-overlay/s6-rc.d/createsuperuser/run
#!/command/with-contenv sh
poetry run -C /opt python manage.py createsuperuser --username=admin --email=  --noinput
