# docker/s6-overlay/s6-rc.d/rabbitmq/run
#!/command/with-contenv sh
rabbitmq-server
